# rpm/voxtype-indicator.service
[Unit]
Description=Voxtype recording indicator (tray icon)
After=voxtype.service graphical-session.target
Requires=graphical-session.target

[Service]
ExecStart=@@LIBEXECDIR@@/voxtype-kde-indicator/voxtype-indicator
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
